# templates/default/redis.upstart.conf.erb
# Redis Datastore Server instance on port <%= port %>
# <%= @name %> generated by Chef
#
description "Upstart job for redis server instance"
author "dennybaa (denis.barishev@gmail.com)"

start on (local-filesystems and runlevel [2345])
stop on runlevel [016]

respawn
respawn limit 10 30

expect fork
limit nofile <%= @ulimit %> <%= @ulimit %>

exec start-stop-daemon --start --chuid <%= @user %>:<%= @group %> --umask 007 --exec <%= File.join(@bin_path, 'redis-server') %> -- <%= @configdir %>/<%= @name %>.conf
